# 99-usbwrite.rules
SUBSYSTEM=="block", ENV{DEVTYPE}=="disk", RUN+="/var/local/usbwriter/udev.monitor"

